# etc/inittab
::sysinit:mkdir /dev/shm
::sysinit:mount -a
::sysinit:/etc/init.d/rcS

::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rcK
::shutdown:umount -a -r

tty1::respawn:/usr/bin/getty 38400 tty1
tty2::respawn:/usr/bin/getty 38400 tty2
tty3::respawn:/usr/bin/getty 38400 tty3
tty4::respawn:/usr/bin/getty 38400 tty4

ttyS0::respawn:/usr/bin/getty 115200 ttyS0
